# data/roles/collectd/templates/upstart.conf.j2
start on runlevel [2345]
stop on runlevel [016]
respawn
exec sh -c '/opt/collectd/sbin/collectd -C /etc/collectd.conf -f >> /var/log/collectd.log 2>&1'
